# inittab
# /etc/inittab
init::sysinit:/sbin/openrc sysinit --quiet
tty1::wait:/bin/bash --rcfile /root/.bashrc +m
bail::ctrlaltdel:/sbin/reboot -f